# var/db/etcupdate/current/etc/rc.d/automountd
#!/bin/sh
#
# $FreeBSD: releng/10.3/etc/rc.d/automountd 289519 2015-10-18 14:59:38Z trasz $
#

# PROVIDE: automountd
# REQUIRE: DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="automountd"
rcvar="autofs_enable"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/${name}"
required_modules="autofs"

load_rc_config $name
run_rc_command "$1"
