# usr/ports/dns/adsuck/files/adsuck.in
#!/bin/sh

# $FreeBSD: tags/RELEASE_10_3_0/dns/adsuck/files/adsuck.in 340872 2014-01-24 00:14:07Z mat $
#
# PROVIDE: adsuck
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable adsuck:
#
# adsuck_enable (bool):   Set to NO by default.

. /etc/rc.subr

name="adsuck"
rcvar=adsuck_enable

command=%%PREFIX%%/sbin/${name}
extra_commands="reload"

load_rc_config $name

: ${adsuck_enable="NO"}
: ${adsuck_flags="-c /var/adsuck -f resolv.conf -r regex Hosts.blc Hosts.pub Hosts.rsk Hosts.trc hosts.small hosts.yoyo"}

run_rc_command "$1"
